# api/samband-api.service
# Samband API systemd service file
# Copy to /etc/systemd/system/samband-api.service
#
# Setup:
#   sudo cp samband-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable samband-api
#   sudo systemctl start samband-api
#
# Logs:
#   sudo journalctl -u samband-api -f

[Unit]
Description=Samband API - Police Events Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/samband-api
Environment="PATH=/opt/samband-api/venv/bin"
ExecStart=/opt/samband-api/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/samband-api/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
